# dib/fluidstack-ironwood/static/etc/systemd/system/ssh.socket.d/ipv6.conf
# Extend ssh.socket to also listen on IPv6.
# The default Debian ssh.socket only has ListenStream=22 (IPv4).
# Ironwood TPU machines have IPv6-only DCN interfaces so SSH is
# unreachable without this override.
[Socket]
ListenStream=
ListenStream=22
ListenStream=[::]:22
